# profile.d/alias.sh
alias ll='ls -la'
alias gitpush='git push --force-if-includes --force-with-lease'
